# tools/tape-network/assets/tape-node.service.tmpl
[Unit]
Description=Tapedrive storage node
After=network.target

[Service]
Type=simple
User=tape
Group=tape
WorkingDirectory={{WORKING_DIR}}
ExecStart={{WORKING_DIR}}/bin/tape-node --config {{WORKING_DIR}}/node.yaml
Restart=on-failure
RestartSec=5s
LimitNOFILE=65536
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
